# tools/wbackup/install/etc/cron.d/wbackup
33 3 * * * wbackup /usr/bin/wbackup || :
