# transmute/install/usr/share/aurex/profile.d/aurex-default-vars.sh
AUREX_RELEASE="10.04"
UBUNTU_CODENAME="lucid"
AUREX_CODENAME="lucid"
UBUNTU_DEFMIRROR="http://archive.ubuntu.com/ubuntu/"
AUREX_DEFMIRROR="http://192.168.100.200/aurex/"
